# root/etc/uci-defaults/40_luci-fb
#!/bin/sh

ver=$(date +%s)
sed -i "s/@ver/$ver/g" /usr/lib/lua/luci/view/fileassistant.htm

rm -f /tmp/luci-indexcache
exit 0
